# packages/tools/bcm2835-bootloader/system.d/rpi3-config-check.service
[Unit]
Description=Fix config.txt if on RPi3

[Service]
ExecStart=/usr/share/bootloader/config-fix.sh

[Install]
WantedBy=multi-user.target
