# packages/network/openssh/system.d/sshd.service
[Unit]
Description=OpenSSH server daemon
ConditionKernelCommandLine=|ssh
After=syslog.target network.target

[Service]
ExecStartPre=/usr/sbin/sshd-keygen
ExecStart=/usr/sbin/sshd -D $SSH_ARGS
ExecReload=/bin/kill -HUP $MAINPID
TimeoutStopSec=1s

[Install]
WantedBy=multi-user.target
