# packages/mediacenter/plexht/system.d/plexht-sources.service
[Unit]
Description=PlexHT sources Setup
Before=plexht.service

[Service]
Type=oneshot
Environment=HOME=/storage
ExecStart=/usr/lib/plexht/plexht-sources
RemainAfterExit=yes

[Install]
WantedBy=plexht.service
